# odoo-longpolling-init.conf
description "odoo-longpolling"
setuid {{ODOO_USER}}

start on startup
respawn
respawn limit 2 5

exec /usr/local/src/odoo/openerp-gevent --config=/etc/odoo/odoo-server.conf
