# configs/odoo-init.conf
description "odoo"
setuid {{ODOO_USER}}

start on startup
respawn
respawn limit 2 5

exec {{ODOO_SOURCE_DIR}}/odoo-bin --config={{OPENERP_SERVER}} --logfile={{LOGS_DIR}}/odoo-server.log
