# configs-docker-host/odoo-docker-db-init.conf
# /etc/init
description "docker odoo db"

start on docker
stop on runlevel [!2345]

respawn
script
  /usr/bin/docker start -a {{DB_ODOO_DOCKER}}
end script