# configs-docker-host/odoo-docker-init.conf
# /etc/init
description "docker odoo"


start on odoo-docker-db
stop on runlevel [!2345]

respawn
script
  /usr/bin/docker start -a {{ODOO_DOCKER}}
end script
Raw