# configs-docker-host/odoo-docker-init.conf
# /etc/init
description "docker runbot"


start on filesystem and started docker-db
stop on runlevel [!2345]

respawn
script
  /usr/bin/docker start -a odoo-runbot
end script
Raw