# configs-docker-host/odoo-docker-db-supervisor.conf
[program:odoo]
command=/usr/local/src/odoo/odoo-bin --config=/etc/odoo/odoo-server.conf --logfile=/var/log/odoo/odoo-server.log
autostart=true
autorestart=true
user={{ODOO_USER}}
directory=/opt/{{ODOO_USER}}/
environment = HOME="/opt/{{ODOO_USER}}/",USER="{{ODOO_USER}}"
