# configs-docker-host/odoo-docker-db-init.conf

# /etc/init
description "docker db"

start on filesystem and started docker
stop on runlevel [!2345]

respawn
script
  /usr/bin/docker start -a db
end script